# deploy/sgac.conf
[program:sgac]
command= /home/toor/sgac/deploy/server.sh
directory=/home/toor/sgac/
autostart=true
autorestart=true
stopasgroup=true
user=toor
stdout_logfile = /home/toor/sgac/logs/supervisor.log